# systemd/agentpeek-tmux.service
[Unit]
Description=agentpeek: dedicated tmux server (own socket + cgroup; survives app restarts)

[Service]
Type=oneshot
RemainAfterExit=yes
# Dedicated socket (-L agentpeek) keeps the server isolated from the user's own
# tmux, and the -f config sets `exit-empty off` so the server stays up with no
# sessions — no placeholder session needed. Because the server lives in THIS
# unit's cgroup, restarting or crashing agentpeek.service never kills it (the
# default KillMode=control-group only reaps the offending unit's own cgroup).
# Stopping THIS unit is the only thing that tears the sessions down.
ExecStart=/usr/bin/tmux -L agentpeek -f @REPO@/conf/agentpeek-server.tmux.conf start-server
ExecStop=/usr/bin/tmux -L agentpeek kill-server

[Install]
WantedBy=default.target
