# launchd/dev.agentpeek.tmux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- agentpeek: dedicated tmux server on its own socket (-L agentpeek).
     Template — setup-macos.sh fills in the placeholders and installs to ~/Library/LaunchAgents.
     macOS counterpart of systemd/agentpeek-tmux.service. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.agentpeek.tmux</string>
  <!-- Start the dedicated server (own socket) as a daemon, then exit (oneshot-
       like). The -f config sets `exit-empty off`, so it stays up with no sessions
       — no placeholder needed. Booting this agent out does NOT kill the server or
       its sessions; the server detaches from launchd. -->
  <key>ProgramArguments</key>
  <array>
    <string>@TMUX@</string>
    <string>-L</string>
    <string>agentpeek</string>
    <string>-f</string>
    <string>@REPO@/conf/agentpeek-server.tmux.conf</string>
    <string>start-server</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>@PATH@</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>@LOGDIR@/tmux.log</string>
  <key>StandardErrorPath</key>
  <string>@LOGDIR@/tmux.log</string>
</dict>
</plist>
